# files/usr/local/etc/rc.d/mackerel_micro
#!/bin/sh
#
# PROVIDE: mackerel_micro
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable mackerel-micro-daemon:
#
#   mackerel_micro_enable="YES"
#   mackerel_micro_conf="/usr/local/etc/mackerel-micro/mackerel-micro.conf"  # optional
#   mackerel_micro_flags=""  # optional extra flags
#

. /etc/rc.subr

name="mackerel_micro"
rcvar="mackerel_micro_enable"
desc="Mackerel Micro Daemon - lightweight Mackerel agent for FreeBSD"

# デフォルト値
: ${mackerel_micro_enable:="NO"}
: ${mackerel_micro_conf:="/usr/local/etc/mackerel-micro/mackerel-micro.conf"}
: ${mackerel_micro_flags:=""}

command="/usr/local/sbin/mackerel-micro-daemon"
command_interpreter="/usr/local/bin/perl"
pidfile="/var/run/${name}.pid"
# ログは logger 経由で syslogd に委譲する。
# newsyslog による世代管理・ローテートは /var/log/messages の既存ルールがそのまま適用される。
# 確認コマンド例: grep 'mackerel-micro' /var/log/messages

start_precmd="mackerel_micro_prestart"
start_cmd="mackerel_micro_start"
stop_cmd="mackerel_micro_stop"
status_cmd="mackerel_micro_status"

mackerel_micro_prestart()
{
    if [ ! -f "${mackerel_micro_conf}" ]; then
        err 1 "Configuration file not found: ${mackerel_micro_conf}" \
              "  Copy the sample: cp ${mackerel_micro_conf}.sample ${mackerel_micro_conf}" \
              "  Then edit it to set MACKEREL_API_KEY."
    fi
}

mackerel_micro_start()
{
    echo "Starting ${name}."
    
    # stdout/stderr の syslog 転送は daemon(8) のネイティブ機能 (-S) に完全委譲する。
    # シェルやパイプを排除することで、PIDファイルに正真正銘の perl のPIDが書き込まれ、
    # 孤児プロセス問題（stopが効かないバグ）を完全に解消する。
    # ※ -o ではなく -S (syslog) を使うため、ローテート時の fd 掴みっぱなし問題も起きない。
    
    /usr/sbin/daemon -f -S -T mackerel-micro -l daemon -s info -p "${pidfile}" \
        /usr/bin/env perl ${command} ${mackerel_micro_conf} ${mackerel_micro_flags}
}

mackerel_micro_stop()
{
    if [ -f "${pidfile}" ]; then
        pid=$(cat "${pidfile}")
        echo "Stopping ${name} (PID: ${pid})."
        kill "${pid}" && rm -f "${pidfile}"
    else
        echo "${name} is not running (no pidfile)."
    fi
}

mackerel_micro_status()
{
    if [ -f "${pidfile}" ] && kill -0 "$(cat ${pidfile})" 2>/dev/null; then
        echo "${name} is running as PID $(cat ${pidfile})."
    else
        echo "${name} is not running."
        return 1
    fi
}

load_rc_config $name
run_rc_command "$1"
